# files/etc/uci-defaults/99-network
uci -q batch << EOI
set network.lan.ipaddr='10.0.0.1'
commit network
EOI
